# daemon/deb/mds-daemon.service
[Unit]
Description=mds-service

[Service]
User=mds
Group=mds
ExecStart=/usr/bin/mds-daemon

[Install]
WantedBy=multi-user.target
